# .config/systemd/user/mpd.service.d/override.conf
[Unit]
After=pulseaudio.service systemd-tmpfiles-setup.service